# run.service
[Unit]
Description = run daemon

[Service]
ExecStart = /root/run.sh


[Install]
WantedBy = multi-user.target