# runzy.service
[Unit]
Description = runzy daemon

[Service]
ExecStart = /root/runzy.sh


[Install]
WantedBy = multi-user.target